# OpenSource-Radio/examples/systemd/votreradio-liquidsoap-restart.service
[Unit]
Description=VotreRadio - Redémarrage quotidien de Liquidsoap

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart liquidsoap
